# deploy/agent-email-layer.service
[Unit]
Description=BuiltByEcho Agent Email Layer API
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=20

[Service]
Type=simple
User=agent-email
Group=agent-email
WorkingDirectory=/opt/agent-email-layer
EnvironmentFile=/etc/agent-email-layer/env
ExecStart=/usr/local/bin/node /opt/agent-email-layer/dist/src/server.js
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
KillSignal=SIGTERM
LimitNOFILE=65536
MemoryMax=512M
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/var/lib/agent-email-layer

[Install]
WantedBy=multi-user.target
