# deploy/agent-email-layer-pg-backup.timer
[Unit]
Description=Nightly Agent Email Layer PostgreSQL backup

[Timer]
OnCalendar=*-*-* 03:15:00 UTC
Persistent=true
RandomizedDelaySec=10m
Unit=agent-email-layer-pg-backup.service

[Install]
WantedBy=timers.target
